# rc.local
## Automatic startup on reboot
# Add a deployment method to /etc/rc.local by sourcing this file like:
#   source rc.local

# Deployment with tmux
deploy="/usr/bin/tmux new-s -d /home/pi/git/piring/ring"
# Alternative deployment without tmux
#deploy="/home/pi/git/piring/ring >>/home/pi/ring.log"

rc=/etc/rc.local
[[ -f $rc ]] &&
	echo "$deploy" |sudo tee -a "$rc" ||
	echo -e "#!/bin/bash\n\n$deploy" |sudo tee "$rc"

# Please check /etc/rc.local to make sure you're starting `ring` only once!
